# config/systemd/asset-observer-medium.service
[Unit]
Description=Asset Observer — Medium Scan (every 3 days, module mapping check)
[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/.openclaw/workspace/scripts/asset_observer.py --depth medium
