# config/systemd/asset-observer-full.service
[Unit]
Description=Asset Observer — Full Scan (every 7 days, snapshot + classification)
[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/.openclaw/workspace/scripts/asset_observer.py --depth full
