# config/systemd/asset-observer-light.timer
[Unit]
Description=Asset Observer — Light Scan (daily)
Requires=asset-observer-light.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
