# config/systemd/baseline-observer.timer
[Unit]
Description=Baseline Observer — Daily

[Timer]
OnCalendar=*-*-* 00:25:00
Persistent=true

[Install]
WantedBy=timers.target
