# config/systemd/baseline-observer.timer
[Unit]
Description=Baseline Observer — Daily
Requires=baseline-observer.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
