# config/systemd/baseline-observer.service
[Unit]
Description=Baseline Observer — Daily baseline drift check
[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/.openclaw/workspace/scripts/baseline_observer.py
