# config/systemd/oek-token-observe.timer
[Unit]
Description=OEK Token Observe (hourly)

[Timer]
OnCalendar=*-*-* *:07:00
Unit=oek-token-observe.service

[Install]
WantedBy=timers.target
