# config/systemd/system-evolution-observer.timer
[Unit]
Description=System Evolution Observer — Daily

[Timer]
OnCalendar=*-*-* 00:35:00
Persistent=true

[Install]
WantedBy=timers.target
