# config/systemd/system-evolution-observer.timer
[Unit]
Description=System Evolution Observer — Daily
Requires=system-evolution-observer.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
